# ops/launchd/com.aegiscore.dogfood.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aegiscore.dogfood</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/york/aegis-core/.venv/bin/aegis-core</string>
        <string>up</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/york/aegis-core</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/york/aegis-core/dogfood/up.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/york/aegis-core/dogfood/up.err.log</string>
</dict>
</plist>
